# files/etc/uci-defaults/99-custom-version
#!/bin/sh
# ============================================================
# ðŸ§© ImmortalWrt å›ºä»¶è‡ªå®šä¹‰ä¿¡æ¯è„šæœ¬
# ä½œè€…: Mr.Zhang
# åŠŸèƒ½:
#  1ï¸âƒ£ è®¾ç½®ä¸»æœºåä¸º Mz.WRT
#  2ï¸âƒ£ LuCI å›ºä»¶ç‰ˆæœ¬æ˜¾ç¤ºè‡ªå®šä¹‰å†…å®¹
#  3ï¸âƒ£ SSH æ¬¢è¿Žè¯­ä¸­æ˜¾ç¤ºç‰ˆæœ¬ä¿¡æ¯ä¸Žæž„å»ºæ—¥æœŸï¼ˆåˆ°æœˆï¼‰
#  4ï¸âƒ£ å¯åŠ¨åŽè‡ªåŠ¨ä¿®å¤ LuCI å›ºä»¶ç‰ˆæœ¬æ˜¾ç¤ºé—®é¢˜
# ============================================================

CUSTOM_HOSTNAME="Mz.WRT"
BUILD_DATE=$(date +%Y.%m)
CUSTOM_VERSION="ImmortalWrt 24.10.4 by Mr.Zhang"

echo "ðŸ§© æ­£åœ¨è®¾ç½®ä¸»æœºåä¸Žç³»ç»Ÿç‰ˆæœ¬ä¿¡æ¯..."

# ------------------ 1ï¸âƒ£ è®¾ç½®ä¸»æœºå ------------------
uci set system.@system[0].hostname="${CUSTOM_HOSTNAME}"
uci commit system
echo "${CUSTOM_HOSTNAME}" > /proc/sys/kernel/hostname

# ------------------ 2ï¸âƒ£ ä¿®æ”¹å›ºä»¶ç‰ˆæœ¬ä¿¡æ¯ ------------------
# /etc/openwrt_release
if [ -f /etc/openwrt_release ]; then
    sed -i "s/DISTRIB_DESCRIPTION='[^']*'/DISTRIB_DESCRIPTION='${CUSTOM_VERSION}'/" /etc/openwrt_release
else
    echo "DISTRIB_DESCRIPTION='${CUSTOM_VERSION}'" >> /etc/openwrt_release
fi

# /etc/os-release
if [ -f /etc/os-release ]; then
    sed -i "s/^PRETTY_NAME=.*/PRETTY_NAME=\"${CUSTOM_VERSION}\"/" /etc/os-release
else
    echo "PRETTY_NAME=\"${CUSTOM_VERSION}\"" >> /etc/os-release
fi

# ------------------ 3ï¸âƒ£ SSH æ¬¢è¿Žè¯­ ------------------
BANNER_FILE="/etc/banner"
{
    echo ""
    echo "============================================"
    echo " ðŸ§© ${CUSTOM_VERSION}"
    echo " ðŸ–¥ï¸  Hostname: ${CUSTOM_HOSTNAME}"
    echo " ðŸ“… Build Date: ${BUILD_DATE}"
    echo "============================================"
} >> "$BANNER_FILE"

# ------------------ 4ï¸âƒ£ é˜²æ­¢ç‰ˆæœ¬ä¿¡æ¯è¢«è¦†ç›– ------------------
RCLOCAL="/etc/rc.local"
grep -q "Mr.Zhang" "$RCLOCAL" 2>/dev/null || cat >> "$RCLOCAL" <<EOF

# ðŸ§© ä¿®å¤ LuCI å›ºä»¶ç‰ˆæœ¬æ˜¾ç¤ºï¼ˆé˜²æ­¢ç³»ç»Ÿè¦†ç›–ï¼‰
if [ -f /etc/openwrt_release ]; then
    sed -i "s/DISTRIB_DESCRIPTION='[^']*'/DISTRIB_DESCRIPTION='${CUSTOM_VERSION}'/" /etc/openwrt_release
fi
if [ -f /etc/os-release ]; then
    sed -i "s/^PRETTY_NAME=.*/PRETTY_NAME=\"${CUSTOM_VERSION}\"/" /etc/os-release
fi
EOF

# ------------------ æ¸…ç†è‡ªèº« ------------------
rm -f /etc/uci-defaults/99-custom.sh
exit 0
